# templates/systemd.service.erb
# This file is managed by Puppet. DO NOT EDIT.
# If you modify this, please also make sure to edit init.sh

[Unit]
Description=InfluxDB is an open-source, distributed, time series database
Documentation=https://docs.influxdata.com/influxdb/
After=<%= @service_after %>

[Service]
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=<%= @service_limit_no_file %>
EnvironmentFile=-<%= @service_defaults %>
ExecStart=/usr/bin/influxd -config <%= @configuration_path %>/<%= @configuration_file %> $INFLUXD_OPTS
KillMode=<%= @service_kill_mode %>
Restart=<%= @service_restart %>

[Install]
WantedBy=<%= @service_wanted_by %>
Alias=<%= @service_alias %>
