# templates/systemd.service.erb
# This file is managed by Puppet. DO NOT EDIT.
# If you modify this, please also make sure to edit init.sh

[Unit]
Description=InfluxDB is an open-source, distributed, time series database
Documentation=https://docs.influxdata.com/influxdb/
After=network-online.target

[Service]
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=65536
EnvironmentFile=-<%= @service_defaults %>
ExecStart=/usr/bin/influxd -config <%= @configuration_path %>/<%= @configuration_file %> $INFLUXD_OPTS
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=influxd.service
